# rootfs/etc/services.d/myservice/run
#!/bin/bash

exec s6-notifyoncheck -c 'sleep 2s' /etc/services.d/myservice/service.sh
