# rootfs/etc/services.d/myservice/run
#!/bin/bash

exec s6-notifyoncheck -c false sleep infinity
